# scripts/crontab
30 1 * * * root /usr/src/app/bin/reporter >> /usr/src/app/logs_report.txt 2>&1
*/10 * * * * root echo 'hello' >> /usr/src/app/hello.txt
